# agentpet-bridge.service
[Unit]
Description=AgentPet Wear OS Bridge
After=network.target

[Service]
Type=simple
# Reemplaza /root/bridge por la ruta real donde clonaste el repo en tu VPS
WorkingDirectory=/root/bridge
ExecStart=/usr/bin/python3 main.py

# ── Configura el agente que quieras usar ──────────────────────────────────
# Hermes Agent (default):
# Environment="AGENT_CMD=hermes chat -Q -q"
# Environment="AGENT_RESUME_CMD=hermes chat --resume {session_id} -Q -q"
#
# llm (simonw/llm):
# Environment="AGENT_CMD=llm -m gpt-4o"
#
# Nombre visible en GET /
# Environment="AGENT_NAME=Hermes"
# API key esperada por el reloj
# Environment="AGENT_API_KEY=tu_clave"
# Alias legacy soportado por compatibilidad
# Environment="HERMES_API_KEY=tu_clave"
# Ajustes del gateway
# Environment="SESSION_MAP_PATH=/root/bridge/watch_sessions.json"
# Environment="WEATHER_LOCATION=Buenos Aires"
# Environment="WATCH_PLATFORM_HINT=This session is coming from a Wear OS watch. Stay fully yourself, but be more direct, concise, expressive, and warm for the watch."
# ─────────────────────────────────────────────────────────────────────────

# Si el CLI del agente no está en el PATH del servicio, añádelo aquí:
# Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/root/.local/bin"

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
